# supervisor.ini
[supervisord]
nodaemon=true

[program:tigervnc]
command=/usr/bin/vncserver :0
autorestart=true

[program:novnc]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900
autorestart=true
